# systemd/cssh.socket
[Unit]
Description=SSH socket

[Socket]
ListenStream=0.0.0.0:24
Accept=yes
