# systemd/csshd@.service
[Unit]
Description=sshd container

[Service]
ExecStart=-/usr/bin/docker run --rm -i frown:local
StandardInput=socket
StandardError=syslog
