# ansible/templates/cssh@.service.j2
[Unit]
Description=sshd container

[Service]
ExecStart=-/usr/bin/docker run --rm --network {{ docker_network_name }} -i {{ chal_container_name }}
StandardInput=socket
StandardError=syslog
