# ansible/templates/cssh.socket.j2
[Unit]
Description=Container SSH Socket

[Socket]
ListenStream={{ container_ssh_port }}
Accept=yes
MaxConnections=1024

[Install]
WantedBy = sockets.target
